# contrib/init/nixd.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env NIXD_BIN="/usr/bin/nixd"
env NIXD_USER="nix"
env NIXD_GROUP="nix"
env NIXD_PIDDIR="/var/run/nixd"
# upstart can't handle variables constructed with other variables
env NIXD_PIDFILE="/var/run/nixd/nixd.pid"
env NIXD_CONFIGFILE="/etc/nix/nix.conf"
env NIXD_DATADIR="/var/lib/nixd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # nixd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$NIXD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run nixd."
        echo "The setting must appear in $NIXD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"NIX Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$NIXD_PIDDIR"
    chmod 0755 "$NIXD_PIDDIR"
    chown $NIXD_USER:$NIXD_GROUP "$NIXD_PIDDIR"
    chown $NIXD_USER:$NIXD_GROUP "$NIXD_CONFIGFILE"
    chmod 0660 "$NIXD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$NIXD_PIDFILE" \
    --chuid $NIXD_USER:$NIXD_GROUP \
    --exec "$NIXD_BIN" \
    -- \
    -pid="$NIXD_PIDFILE" \
    -conf="$NIXD_CONFIGFILE" \
    -datadir="$NIXD_DATADIR" \
    -disablewallet \
    -daemon

